# dashboard/network-pulse.service
[Unit]
Description=Network Pulse Dashboard - Live Device Monitor
After=network.target

[Service]
Type=simple
User=yb
Group=yb
WorkingDirectory=/usr/local/yb/.Network/dashboard
ExecStart=/usr/bin/python3 /usr/local/yb/.Network/dashboard/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
